# config/goweb.service
[Unit]
Description=goweb
[Service]
Type=simple
Restart=always
RestartSec=5s
ExecStart=/go/go-web/main
[Install]
WantedBy=multi-user.target